# launchd/com.dannytsao.personalkm.copilot-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.dannytsao.personalkm.copilot-ingest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/dannytsao/Library/Application Support/PersonalKM/run_mac_mini_copilot_ingest.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/dannytsao/Documents/PersonalKM/Personalkm-vault</string>

  <!-- Daily at 03:00. If the Mac is off at that time, launchd does NOT
       catch up the missed run on next boot — it simply waits for the
       next 03:00 (unlike StartInterval, which would fire immediately to
       make up elapsed time). RunAtLoad is intentionally omitted so
       loading/reloading this job never triggers an extra run. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/dannytsao/Library/Logs/PersonalKM/copilot-ingest.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/dannytsao/Library/Logs/PersonalKM/copilot-ingest.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PERSONALKM_REPO_ROOT</key>
    <string>/Users/dannytsao/Documents/GitHub/DannyTsao/PersonalKM</string>
    <key>PERSONALKM_VAULT_ROOT</key>
    <string>/Users/dannytsao/Documents/PersonalKM/Personalkm-vault</string>
  </dict>
</dict>
</plist>
